# private_dot_config/systemd/user/modprobed-db-recall.timer
[Unit]
Description=Run modprobed-db recall on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
