# private_dot_config/systemd/user/modprobed-db-recall.service
[Unit]
Description=loads used modules on boot
[Service]
Type=oneshot
ExecStart=sudo modprobed-db recall
