# private_dot_config/systemd/user/timers.target.wants/borg.timer
[Unit]
Description=Run borg backup

[Timer]
#OnCalendar=daily
#Persist=true
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
